# iso-builder/config/includes.chroot/etc/systemd/system/installer.target
[Unit]
Description=DebDox Installer Target
Requires=basic.target
After=basic.target
AllowIsolate=yes
